# supervisor/webapp.conf
; =============================================================================
; CivicAssistBot v5.0 — Supervisor Configuration
; =============================================================================
;
; ЄДИНИЙ ПРОЦЕС: bot.py запускає polling, HTTP-сервер (порт 4243),
; Stripe webhook і планувальник Termin — все в одному.
;
; ВСТАНОВЛЕННЯ:
;   1. Відредагуйте APP_DIR та USER нижче під своє середовище
;   2. cp supervisor/webapp.conf /etc/supervisor/conf.d/civicassistbot.conf
;   3. supervisorctl reread && supervisorctl update
;   4. supervisorctl start civicassistbot
;
; =============================================================================

[program:civicassistbot]
; --- Команда ---
; Використовуємо python з venv якщо він є, інакше системний python3.
; Змініть шлях якщо venv розташований в іншому місці.
command=/opt/civicassistbot/venv/bin/python bot.py

; --- Робоча директорія (ЗМІНІТЬ) ---
directory=/opt/civicassistbot

; --- Користувач (ЗМІНІТЬ, не використовуйте root у production) ---
user=botuser

; --- Автозапуск ---
autostart=true
autorestart=true
; Вважати процес запущеним через 10 секунд після старту
startsecs=10
; Спроб перезапуску перед переходом в FATAL
startretries=5

; --- Graceful shutdown ---
; Дати 30 секунд на завершення активних Stripe-запитів і Termin-моніторингу
stopwaitsecs=30
stopsignal=TERM
killasgroup=true
stopasgroup=true

; --- Логування ---
stdout_logfile=/opt/civicassistbot/logs/bot_stdout.log
stdout_logfile_maxbytes=20MB
stdout_logfile_backups=5
stderr_logfile=/opt/civicassistbot/logs/bot_stderr.log
stderr_logfile_maxbytes=20MB
stderr_logfile_backups=5

; --- Змінні середовища ---
; Секрети зберігайте в .env, не тут.
; Supervisor не завантажує .env автоматично — використовуйте EnvFile у systemd
; або обгортку start_bot.sh (див. нижче).
environment=
    PYTHONUNBUFFERED="1",
    PYTHONPATH="/opt/civicassistbot"

; --- Пріоритет ---
priority=100
